# rpm/netrek.crontab
#!/bin/bash

confirm=`grep CONFIRM /usr/games/netrek/.sysdef`

#
# Should ALWAYS be in there, but check just in case
#
if [ $? -ne 0 ]; then
  exit 0
fi

#
# What is the value of CONFIRM?
#
num=`echo $confirm | cut -d '=' -f2`
if [ $? -ne 0 ]; then
  exit 0
fi

if [ $num -gt 0 ]; then
  su -c "/usr/games/netrek/tools/updatekeys" - games
fi

exit 0
